# tools/appserver/qdsapp.socket.in
[Unit]
Description=QtDataSync AppServer Socket
Documentation=https://github.com/Skycoder42/QtDataSync
After=network.target
PartOf=$${APPNAME}.service

[Socket]
ListenStream=4242

[Install]
WantedBy=sockets.target
